# config/conf/proc_2.conf
[program:proc_2]
priority=3
command=/bin/bash -c 'python3 /opt/womginx/config/self-ping.py'
priority=995